# backup/systemd/qdrant-backup.timer
[Unit]
Description=Hourly Qdrant snapshot + S3 backup

[Timer]
OnCalendar=hourly
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
